# supervisor/supervisor.apps.ini
[program:socketio]
directory = /track_a_tron
command = hypercorn app:create_app --config hypercorn/config.toml
autostart = false
autorestart = false
startretries = 0

[program:huey]
directory = /track_a_tron
command = huey_consumer app.huey.run -w 3 --logfile=/track_a_tron/logs/huey.log
autostart = false
autorestart = false
startretries = 0