# supervisor/supervisor.apps.ini
[program:gunicorn]
directory = /track_a_tron
command = gunicorn --config gunicorn/gunicorn.conf.py
autostart = false
autorestart = false
startretries = 0
